# ksystemstats/plasma-ksystemstats.service.in
[Unit]
Description=Track hardware statistics
PartOf=graphical-session.target

[Service]
ExecStart=@KDE_INSTALL_FULL_BINDIR@/ksystemstats
BusName=org.kde.ksystemstats
Slice=background.slice
